# crontab
* 15 * * * python /path/to/repo/amazon_auth.py
